# deploy/orbit-os.service
[Unit]
Description=Orbit OS — PI and Codex VPS bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=codex
Group=codex
WorkingDirectory=/home/codex/agentic-os
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=ORBIT_DATA_DIR=/var/lib/orbit-os
EnvironmentFile=/etc/agentic-os/orbit.env
ExecStart=/usr/local/bin/node /home/codex/agentic-os/server.mjs
Restart=always
RestartSec=3
TimeoutStopSec=20
UMask=0077
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/.pi /var/lib/orbit-os /home/codex/agentic-os
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# NoNewPrivileges cannot be enabled yet: the PI/Codex bridge deliberately uses
# narrowly allowlisted sudo commands to create isolated transient units.

[Install]
WantedBy=multi-user.target
