# deploy/orbit-os.service
[Unit]
Description=Orbit OS — PI and Codex VPS bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=codex
Group=codex
WorkingDirectory=/home/codex/agentic-os
Environment=NODE_ENV=production
EnvironmentFile=/etc/agentic-os/orbit.env
ExecStart=/usr/bin/node /home/codex/agentic-os/server.mjs
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
